# evjs.rules
ACTION=="add", KERNEL=="event*", SUBSYSTEM=="input", RUN+="/usr/bin/evjscal -c '%E{DEVNAME}'"
